# SUSE/podman-build.service
[Unit]
Description=Podman Build Service

[Service]
ExecStart=$HOME/Developer/Distrobox/SUSE/podman_build.sh
Type=oneshot
