# Fedora/podman-build.service
[Unit]
Description=Podman Build Service

[Service]
ExecStart=/home/brian/Developer/DistroboxImages/Fedora/podman_build.sh
Type=oneshot
