# Fedora/podman-build.timer
[Unit]
Description=Run Podman Build on Boot

[Timer]
OnBootSec=5min

[Install]
WantedBy=default.target
